# deploy/systemd/clickhouse-spool-flush.service
[Unit]
Description=Flush ClickHouse Spooled Writes
After=network.target clickhouse-server.service
Wants=clickhouse-server.service

[Service]
Type=oneshot
User=uptrack
Group=uptrack
Environment=CLICKHOUSE_HOST=100.C.C.C
Environment=CLICKHOUSE_PORT=8123
ExecStart=/usr/local/bin/clickhouse-flush-spool.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=clickhouse-spool-flush
